# deploy/wx01.service
[Unit]
Description=wx01 weather station collector
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
ExecStartPre=/usr/local/bin/wx01-db start
ExecStart=/usr/local/bin/wx01
AmbientCapabilities=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5
Environment=WX01_DATABASE_URL=postgres://wx01:wx01@127.0.0.1:5433/wx01?sslmode=disable
Environment=WX01_UDP_PORT=50222
Environment=WX01_HTTP_PORT=80

[Install]
WantedBy=multi-user.target
